# systemd/oneplace-youtube-analytics.timer
[Unit]
Description=OnePlace YouTube analytics (daily)

[Timer]
OnCalendar=*-*-* 07:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
