# systemd/oneplace-podcast-publish.service
[Unit]
Description=OnePlace podcast publish (R2 RSS)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/OnePlace
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/OnePlace/.venv/bin/python -m app.worker publish-podcast --all-projects
